# homelab-dns/com.tmilata.homelab-dns.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tmilata.homelab-dns</string>

    <!-- Run the long-running watcher. It blocks on `route -n monitor` and only
         re-evaluates DNS when the network path actually changes — no polling,
         no file-watching, idle (zero CPU) on a stable network. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/sbin/homelab-dns-watch.sh</string>
    </array>

    <!-- Start at boot and keep it alive (restart if it ever exits/crashes). -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- If it does exit, don't respawn faster than this. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
